# jenkins/config/supervisor.ini
[program:filebeat]
command=/filebeat/filebeat-7.6.2-linux-x86_64/filebeat -c /filebeat/filebeat-7.6.2-linux-x86_64/filebeat.yml

[program:jenkins]
command=/usr/local/bin/jenkins.sh
user=jenkins
redirect_stderr=true
stdout_logfile=/var/log/supervisor/jenkins-supervisor-stdout.log